# admin-scripts/build/rpms/rbuild-ids/SOURCES/rbuild-ids/etc/cron.daily/run-ids
#!/bin/bash
#
# run-ids - Runs intrusion detection system from cron job
#

/usr/local/bin/check-ids.sh
